# systemd/sloc.socket
[Unit]
Description=FAUST CTF SL Online Compiler Service

[Socket]
ListenStream=64646
Accept=true

[Install]
WantedBy=sockets.target

